# crontab
* * * * * python3 /home/username/clean-tokens.py
